# roles/otel/templates/otel-collector.service.j2
[Unit]
Description=OpenTelemetry Collector
After=network.target

[Service]
Type=simple
User=otel-collector
Group=otel-collector
ExecStart=/usr/bin/otel-collector --config=/etc/otel-collector/config.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
